# ops/planning-backup.service
[Unit]
Description=Planning Previsionnel - Daily PostgreSQL Backup
After=network.target

[Service]
Type=oneshot
EnvironmentFile=/opt/planning_previsionnel/shared/api.env
ExecStart=/opt/planning_previsionnel/current/ops/backup_postgres_daily.sh
User=postgres
Group=postgres
